# userspace/udev/01-forensic-readonly.rules
# Mark new block devices as read-only
ACTION=="add", SUBSYSTEM=="block", KERNEL!="ram*", RUN+="/usr/sbin/wrtblk %k"
